# weblogic/init/server.conf
start on runlevel [345]

exec /bin/su - oracle -- @DOMAIN_HOME/bin/startManagedWebLogic.sh "@SERVER"
